# cli-proxy-api-plus.service
[Unit]
Description=CLIProxyAPI Plus - Proxy API Server for AI CLI Models
Documentation=https://github.com/router-for-me/CLIProxyAPI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cliproxy
Group=cliproxy
WorkingDirectory=/opt/cli-proxy-api-plus

# Main service binary
ExecStart=/opt/cli-proxy-api-plus/cli-proxy-api-plus

# Restart policy - restart on failure
Restart=on-failure
RestartSec=5s
StartLimitInterval=300
StartLimitBurst=5

# Environment file for configuration
EnvironmentFile=-/etc/cli-proxy-api-plus/environment
EnvironmentFile=-/opt/cli-proxy-api-plus/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cli-proxy-api-plus /var/log/cli-proxy-api-plus
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
TasksMax=512

# Memory and CPU limits (adjust based on your needs)
MemoryMax=2G
MemoryHigh=1.5G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cli-proxy-api-plus

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
